# meta-openeuler/recipes-core/os-base/os-base/inittab
::sysinit:/etc/rc.d/rc.sysinit
::respawn:/sbin/getty -L 115200 ttyAMA0 vt102
::ctrlaltdel:/bin/false
