# meta-openeuler/recipes-core/os-base/os-base/rc.local
#!/bin/bash

/sbin/auditd
/sbin/syslogd
