# meta-openeuler/recipes-core/os-base/os-base/rc.local
#!/bin/bash

. /etc/init.d/auditd start
/usr/sbin/rsyslogd

chmod 640 /var/log/*
chmod 640 /dev/mem
chmod 600 /var/log/audit/audit.log

. /etc/init.d/sshd start 1>/dev/null 2>&1 &
